# os customizing/etc/profile.d/motd.sh
#!/bin/bash
#
echo -e "
#########################################
#
# i-Vertix 2017 Network Management System
# 2.4.0
#
# All rights reserved
# (c)2016 PGUM GmbH/Srl
# info@pgum.eu, support@pgum.eu
# www.pgum.eu
# 
# Welcome to `hostname`
#
# This system is running `cat /etc/redhat-release`
# kernel is `uname -r`
#
# You are logged in as `whoami`
#
#########################################
"