# os customizing/etc/profile.d/motd.sh
#!/bin/bash
#
echo -e "
##################################
#
# Osiris 2.2.0 NMS Poller
#
# All rights reserved
# (c)2014, BI-S des Jochen Platzgummer
# info@bi-s.it, support@bi-s.it
# www.bi-s.it
# 
# Welcome to `hostname`
#
# This system is running `cat /etc/redhat-release`
# kernel is `uname -r`
#
# You are logged in as `whoami`
#
##################################
"